# infra/systemd/zetta-wallet-rollup.timer
[Unit]
Description=Refresh Zetta wallet trade rollup every minute

[Timer]
OnBootSec=3min
OnUnitInactiveSec=1min
AccuracySec=10s
Unit=zetta-wallet-rollup.service

[Install]
WantedBy=timers.target
