# infra/systemd/zetta-load.timer
[Unit]
Description=Run capped Zetta raw loaders every 15 minutes

[Timer]
OnBootSec=2min
OnUnitInactiveSec=15min
AccuracySec=30s
Unit=zetta-load.service

[Install]
WantedBy=timers.target
