# infra/systemd/zetta-wallet-rollup.service
[Unit]
Description=Zetta realtime wallet trade rollup
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner wallet-rollup
MemoryAccounting=true
MemoryMax=8G
TasksMax=64
TimeoutStartSec=5min
