# infra/systemd/zetta-chain-frontier.service
[Unit]
Description=Zetta realtime Polymarket chain log frontier
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner chain-frontier
MemoryAccounting=true
MemoryMax=4G
TasksMax=64
TimeoutStartSec=5min
