# infra/systemd/zetta-marts.service
[Unit]
Description=Zetta analytical mart builder
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner marts
MemoryAccounting=true
MemoryMax=24G
TasksMax=128
TimeoutStartSec=30min
