# infra/systemd/zetta-wallet-pnl-candidates.timer
[Unit]
Description=Seed Zetta smart wallet PnL candidates every minute

[Timer]
OnBootSec=2min
OnUnitInactiveSec=1min
AccuracySec=15s
Unit=zetta-wallet-pnl-candidates.service

[Install]
WantedBy=timers.target
