# infra/systemd/zetta-wallet-pnl-candidates.service
[Unit]
Description=Zetta smart wallet candidate PnL seeder
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner wallet-pnl-candidates
MemoryAccounting=true
MemoryMax=1G
TasksMax=64
TimeoutStartSec=5min
