# infra/systemd/zetta-frontier-trades.service
[Unit]
Description=Zetta realtime trade task seeder
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner frontier-trades
MemoryAccounting=true
MemoryMax=1G
TasksMax=64
TimeoutStartSec=2min
