# infra/systemd/zetta-load-trades-realtime.service
[Unit]
Description=Zetta realtime trade loader
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner load-trades-realtime
MemoryAccounting=true
MemoryMax=8G
TasksMax=64
TimeoutStartSec=5min
