# infra/systemd/zetta-chain-pnl.service
[Unit]
Description=Zetta chain-verified wallet PnL mart
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner chain-pnl
MemoryAccounting=true
MemoryMax=8G
TasksMax=64
TimeoutStartSec=10min
