# infra/systemd/zetta-wallet-screener.service
[Unit]
Description=Zetta all-wallet whale and smart wallet screener
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner wallet-screener
MemoryAccounting=true
MemoryMax=12G
TasksMax=64
TimeoutStartSec=30min
