# infra/systemd/zetta-frontier.service
[Unit]
Description=Zetta recent event task seeder
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner frontier
