# infra/systemd/zetta-wallet-pnl-load.service
[Unit]
Description=Zetta wallet PnL raw loader
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner wallet-pnl-load
MemoryAccounting=true
MemoryMax=2G
TasksMax=64
TimeoutStartSec=10min
