# infra/systemd/zetta-wallet-refresh.service
[Unit]
Description=Zetta realtime wallet refresh task seeder
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner wallets
MemoryAccounting=true
MemoryMax=1G
TasksMax=64
TimeoutStartSec=2min
