# infra/systemd/zetta-ws-market.service
[Unit]
Description=Zetta CLOB market WebSocket raw collector
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/zetta/zetta.env
ExecStart=/usr/local/bin/zetta-runner ws-market
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
